# super_visor.ini
[program:#script_name#]   ;脚本名称
command=#startup_script_command# ;执行的启动指令
directory=#script_work_dir#   ;脚本工作目录
autorstart=true     ;supervisor启动后，自动启动该脚本
autorestart=true    ;脚本进程挂了后，自动拉起

startretries=3      ;脚本进程挂了后，自动拉起重试次数

stopwaitsecs=10     ;发送结束信号后，最大等待时长，一但超过，强制KILL

stdout_logfile="#stdout_logfile#"    ;脚本输出 重定向
stderr_logfile="#stderr_logfile#"    ;脚本错误 重定向
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB

user=root
numprocs=1

process_name=#process_name# ;启动后，进程名称